# templates/cimetrics.conf.erb
description "CIMetrics Aggregator Service"
author "Gabriel Loewen <t-gabloe@microsoft.com"

start on runlevel [234]
stop on runlevel [0156]

exec /CIMetricsTool/statsproj/openstack_stats/aggregatorService.py

respawn
